# source/sec/sec-daemon/sec-daemon.service
[Unit]
Description=SEC daemon Service Script

[Service]
User=root
ExecStart='/usr/bin/secd'

[Install]
WantedBy=multi-user.target

